# systemd/fischverkauf.service
[Unit]
Description=Fischerfest Kassensystem
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/fischverkauf
ExecStart=/home/pi/fischverkauf/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
